# novnc/conf.d/icewm.conf
[program:icewm]
command=icewm
autorestart=true
